# quadlets/jellyfin-sidecar.container
[Unit]
Description=Jellyfin Sidecar - Prevents shutdown during streaming
After=jellyfin.service
BindsTo=jellyfin.service

[Container]
Image=ghcr.io/addisonbair/homelab-sidecars:jellyfin
ContainerName=jellyfin-sidecar
Network=host
Environment=JELLYFIN_URL=http://localhost:8096
Environment=JELLYFIN_API_KEY_FILE=/secrets/jellyfin-api-key
Environment=JELLYFIN_GRACE_PERIOD=5m
Environment=POLL_INTERVAL=30s
Environment=INHIBIT_WHAT=shutdown:sleep
Volume=/run/dbus:/var/run/dbus:ro
Volume=/etc/homelab:/secrets:ro,z
SecurityLabelDisable=true

[Service]
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
