# deploy/raid-inhibitor.service
[Unit]
Description=RAID Health Inhibitor
Documentation=https://github.com/addisonbair/homelab-sidecars
After=local-fs.target
Wants=local-fs.target

[Service]
Type=simple
ExecStart=/usr/local/bin/raid-inhibitor -arrays=md0 -interval=60s
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/proc/mdstat

[Install]
WantedBy=multi-user.target
